# packaging/systemd/gic-ipsec-helper.service.example
[Unit]
Description=GIC IPsec helper example

[Service]
Type=oneshot
ExecStart=/usr/libexec/gic-ipsec-client/gic-ipsec-helper --help
