# packaging/systemd/gic-ipsec-helper.service.example
[Unit]
Description=GIC IPsec helper example
Documentation=man:swanctl(8)

[Service]
Type=oneshot
ExecStart=/usr/bin/gic-ipsec-helper --help

# The MVP uses pkexec for helper invocation. This unit is an example placeholder
# for future distros that prefer a narrow systemd-mediated helper flow.
