# packaging/systemd/see-ipsec-helper.service.example
[Unit]
Description=SEE IPsec helper example

[Service]
Type=oneshot
ExecStart=/usr/libexec/see-ipsec-client/see-ipsec-helper --help
